# configs/audioserver/systemd/netsim_moos_gateway-arctic.service
[Unit]
Description=NetSim Goby to/from MOOS gateway (Arctic Environment)
After=gobyd.service
BindsTo=gobyd.service

[Service]
ExecStart=/opt/goby3/build/bin/goby_moos_gateway /opt/modemsim/missions/goby_moos_gateway-arctic.pb.cfg 

[Install]
WantedBy=gobyd.service

